# apt-monitor_1.0-1_all/DEBIAN/postinst
#!/bin/bash
systemctl daemon-reload
systemctl enable apt-monitor
systemctl restart apt-monitor
