# apt-monitor_1.0-1_all/DEBIAN/prerm
#!/bin/bash
# Para o serviço e desabilita a inicialização junto com o sistema
systemctl stop apt-monitor || true
systemctl disable apt-monitor || true
systemctl daemon-reload || true
